# deploy/veloseller-web.service
[Unit]
Description=Veloseller Web (Next.js)
After=network.target veloseller-worker.service

[Service]
Type=simple
User=veloseller
Group=veloseller
# Sandbox (defense-in-depth). NoNewPrivileges блокирует эскалацию через setuid/sudo,
# обрубая путь «RCE в Next.js → sudo finalize.sh → root». ФС только на чтение, кроме
# /opt/veloseller (там .next/cache); приватный /tmp.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/veloseller
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
WorkingDirectory=/opt/veloseller/apps/web
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/veloseller/apps/web/.env.production
# next-бинарь лежит в hoisted node_modules корня монорепо (npm workspaces)
ExecStart=/usr/bin/node /opt/veloseller/node_modules/next/dist/bin/next start -p 3000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
