# deploy/veloseller-web@.service
[Unit]
Description=Veloseller Web (Next.js) :%i
After=network.target veloseller-worker.service

[Service]
Type=simple
User=veloseller
Group=veloseller
# Sandbox (defense-in-depth). NoNewPrivileges блокирует эскалацию через setuid/sudo,
# обрубая путь «RCE в Next.js → sudo finalize.sh → root». ФС только на чтение, кроме
# /opt/veloseller (там .next/cache); приватный /tmp.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/veloseller
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
WorkingDirectory=/opt/veloseller/apps/web
EnvironmentFile=/opt/veloseller/apps/web/.env.production
Environment=NODE_ENV=production
Environment=PORT=%i
# next-бинарь в hoisted node_modules корня монорепо (npm workspaces).
# Имя инстанса = порт: veloseller-web@3001 → слушает :3001.
ExecStart=/usr/bin/node /opt/veloseller/node_modules/next/dist/bin/next start -p %i
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
